# .config/systemd/user/xdg-desktop-portal-gnome.service.d/override.conf
[Service]
Environment="XDPH_FORCE_SCREENCAST_CHOOSER=0"
